# systemd/meshbot.service
[Unit]
Description=Groq AI Mesh Chatbot
After=network.target

[Service]
Environment=PYTHONPATH=/home/YOUR_USER/Meshtastic-PiAgent
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/YOUR_USER/Meshtastic-PiAgent/mesh_groq_ai_bot_oled.py
WorkingDirectory=/home/YOUR_USER/Meshtastic-PiAgent
StandardOutput=append:/home/YOUR_USER/Meshtastic-PiAgent/meshbot.log
StandardError=append:/home/YOUR_USER/Meshtastic-PiAgent/meshbot.log
Restart=always
RestartSec=5
StartLimitIntervalSec=0
User=YOUR_USER

[Install]
WantedBy=multi-user.target
