# systemd/mode-selector.service
[Unit]
Description=Pi-Bot Boot Mode Selector (Waveshare 1.44" LCD HAT)
After=local-fs.target sysinit.target
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/usr/bin/python3 /root/mode_selector.py
User=root
Environment=PYTHONPATH=/root/Raspyjack
Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
